# jrp-client/src/bin/jrp-client.service
[Unit]
Description=JRP Client Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/jrp-client
ExecStart=/usr/bin/java -Dfile.encoding=utf-8 -Dspring.config.location=./application.yml -jar jrp-client-1.0.2.jar
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target